# service/pgbouncer/finish
#!/bin/sh

exec 2>&1
rm -f /run/postgresql/pgbouncer.run
