# service/pg_auto_failover/finish
#!/bin/sh

exec 2>&1
rm -f /run/postgresql/pg_auto_failover.run
