# service/postgres/finish
#!/bin/sh

exec 2>&1
rm -f /run/postgresql/postgres.run
