# service/postgres/run
#!/bin/sh

exec 2>&1
if [ -z "$PG_AUTO_FAILOVER" ]; then
    install -d -m 0750 -o "$USER" -g "$GROUP" "$PGDATA"
else
    install -d -m 1775 -o "$USER" -g "$GROUP" "/tmp/pg_autoctl$PGDATA"
fi
install -d -m 1775 -o "$USER" -g "$GROUP" /run/postgresql /var/log/postgresql
if [ "$(pgrep runsvdir | wc -l)" -gt "0" ]; then
    install -d -m 0750 -o "$USER" -g "$GROUP" "$HOME/pg_arclog" "$HOME/pg_log" "$HOME/pg_rman"
    chmod 755 supervise
    chown "$USER":"$GROUP" supervise/ok supervise/control supervise/status
fi
rm -f /run/postgresql/postgres.run "$PGDATA/postmaster.pid"
realpath "$0"
set -ex
test -d "$PGDATA/base" || chpst -u "$USER":"$GROUP" /etc/service/postgres/init
test -d "$PGDATA/base"
if [ -z "$PG_AUTO_FAILOVER" ]; then
    rm -rf /etc/service/pglisten
    exec chpst -u "$USER":"$GROUP" -L /run/postgresql/postgres.run postmaster
else
    if [ -z "$PG_AUTO_FAILOVER_MONITOR" ]; then
        rm -rf /etc/service/pgbouncer
        rm -rf /etc/service/pglisten
    fi
    exec chpst -u "$USER":"$GROUP" -L /run/postgresql/postgres.run pg_autoctl run
fi
