# service/postgres/run
#!/bin/sh

exec 2>&1
#install -d -m 0750 -o "$USER" -g "$GROUP" "$PGDATA"
install -d -m 1775 -o "$USER" -g "$GROUP" /run/postgresql /var/log/postgresql
#realpath "$0"
#set -ex
chmod 755 supervise
chown "$USER":"$GROUP" supervise/ok supervise/control supervise/status
#exit
#rm -f "$PGDATA/postmaster.pid"
realpath "$0"
set -ex
if [ -z "$PG_AUTO_FAILOVER_MONITOR" ]; then
    chpst -u "$USER":"$GROUP" pg_autoctl -vvv create monitor --nodename "$(hostname)" --no-ssl --skip-pg-hba
    chpst -u "$USER":"$GROUP" /etc/service/postgres/conf
    chpst -u "$USER":"$GROUP" /etc/service/postgres/hba
elif [ ! -d "$PGDATA/base" ]; then
    chpst -u "$USER":"$GROUP" pg_autoctl -vvv create postgres --nodename "$(hostname)" --no-ssl --skip-pg-hba --monitor="$PG_AUTO_FAILOVER_MONITOR"
    if PGTARGETSESSIONATTRS=read-write chpst -u "$USER":"$GROUP" psql -qtc "select now()"; then
        chpst -u "$USER":"$GROUP" createuser pgautofailover_monitor
        chpst -u "$USER":"$GROUP" pg_ctl --wait --mode=fast stop
        chpst -u "$USER":"$GROUP" /etc/service/postgres/conf
        chpst -u "$USER":"$GROUP" /etc/service/postgres/hba
    else
        chpst -u "$USER":"$GROUP" pg_ctl --wait --mode=fast stop
    fi
fi
exec chpst -u "$USER":"$GROUP" pg_autoctl -vvv run
